# homedir/.config/systemd/user/gnome-keyring-secrets.service
# Complete the PAM-started login keyring from the sway session.
#
# pam_gnome_keyring (auto_start, in the greetd login stack) launches the
# "--login" daemon and unlocks it with the login password. But the step that
# registers the Secret Service on the bus -- "gnome-keyring-daemon --start" --
# normally comes from GNOME's autostart, whose .desktop files are gated
# OnlyShowIn=GNOME;Unity;MATE and so never run under sway. Without it, the first
# app to request org.freedesktop.secrets makes D-Bus spawn a *second* daemon
# (the one that logs "discover_other_daemon"). Running --start here connects to
# the existing login daemon instead of spawning a new one, so exactly one daemon
# owns the bus name.
#
# ssh is deliberately omitted (keychain/gpg-agent own SSH_AUTH_SOCK); pkcs11 and
# secrets are the components apps actually consume.
[Unit]
Description=Register the login gnome-keyring on the session bus
PartOf=graphical-session.target
After=graphical-session.target
# The Secret Service portal impl activates this name; claim it before apps ask.
Before=org.freedesktop.secrets.service

[Service]
# --start connects to the already-running login daemon and returns immediately
# (it does not fork a resident process), so this is a oneshot that just needs to
# have run before apps request secrets.
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/gnome-keyring-daemon --start --components=pkcs11,secrets

[Install]
WantedBy=graphical-session.target
