# homedir/.config/systemd/user/xdg-desktop-portal.service.d/override.conf
[Service]
Environment="XDG_SESSION_TYPE=wayland"
Environment="XDG_SESSION_DESKTOP=sway"
Environment="XDG_CURRENT_DESKTOP=sway"